# supervisor/conf.d/ss-manager.conf
[program:ss-manager]
command=/usr/bin/ss-manager -m aes-256-cfb -u --manager-address 127.0.0.1:6001 -v
user=root
stdout_logfile=/var/log/supervisor/ss-manager.log
redirect_stderr=true
priority=101

