# supervisor/conf.d/ss-server.conf
[program:ss-server]
command=/usr/bin/ss-server -c /etc/shadowsocks/ss-server.json -v 
user=root
stdout_logfile=/var/log/supervisor/ss-server.log
redirect_stderr=true
priority=100

