# supervisor/conf.d/ssmgr-webui.conf
[program:ssmgr-webui] 
command=/usr/bin/ssmgr -c /etc/shadowsocks/webui.yml --debug
user=root
stdout_logfile=/var/log/supervisor/webui.log
redirect_stderr=true
priority=103

