# supervisor/conf.d/ssmgr.conf
[program:ssmgr] 
command=/var/www/shadowsocks-manager/bin/ssmgr -c /etc/shadowsocks/ssmgr.yml --debug
user=root
stdout_logfile=/var/log/supervisor/ssmgr.log
redirect_stderr=true
priority=102
